# docker/novnc/conf.d/i3.conf
[program:i3]
command=i3
autorestart=true