# novnc/conf.wm/dwm.conf
[program:dwm]
command=dwm
autorestart=true
