# novnc/conf.wm/xfce.conf
[program:xfce]
command=startxfce4
autorestart=true